# deploy/ict-telegram-bot.service
# ============================================
# ICT TELEGRAM BOT SYSTEMD SERVICE
# Telegram UI process - runs separately from main engine
# Save to: /etc/systemd/system/ict-telegram-bot.service
# ============================================

[Unit]
Description=ICT Trading Bot - Telegram UI
After=network.target ict-bot.service
Wants=network-online.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/YOUR_REPO
ExecStart=/home/ubuntu/YOUR_REPO/.venv/bin/python -m src.telegram_bot
Restart=always
RestartSec=15
StandardOutput=append:/home/ubuntu/YOUR_REPO/runtime_logs/telegram.log
StandardError=append:/home/ubuntu/YOUR_REPO/runtime_logs/telegram_error.log
EnvironmentFile=/home/ubuntu/YOUR_REPO/.env

[Install]
WantedBy=multi-user.target
