# deploy/ict-exchange-fills-pull.service
# ============================================
# ICT EXCHANGE-FILLS PULL (BL-20260713-EXCHANGE-FILLS-STORE-EMPTY)
# Refreshes runtime_state/exchange_fills.sqlite with the last 7 days of REAL
# Bybit fills for bybit_2 (category=linear — USDT-margined perps) so the
# exchange-truth P&L surface /api/bot/pnl/exchange stays populated instead of
# serving zeros.
#
# Why this exists: the S-067 puller (scripts/pull_exchange_fills.py) shipped
# with NO timer and NO system-action, AND pulled the spot category only — so
# the fills store had never accrued a single row while bybit_2 traded linear
# perps daily (the 2026-07-13 health-review found /api/bot/pnl/exchange dark
# with 0 fills/7d). PR #6290 added --category linear + the on-demand
# pull-exchange-fills action; this unit makes the refresh durable so the
# surface doesn't go stale again the moment the on-demand run ages out of the
# 7-day window.
#
# LIVE-TRADER BOX ONLY. Reads Bybit over the public API using BYBIT_API_KEY_2 /
# BYBIT_API_SECRET_2 from the trader .env; read-only on the exchange side
# (fetch_my_trades). Touches NO service and NO trade_journal.db table. The
# store's primary key is exec_id, so the daily 7-day over-sampling is
# idempotent (duplicate fills are dropped).
#
# DATA_DIR: this unit MUST resolve runtime_state_dir() to the SAME root as
# ict-web-api (the reader) — /data/bot-data. Set explicitly below (mirrors
# ict-mes-ibkr-pull.service) so the writer and the API agree on
# $DATA_DIR/runtime_state/exchange_fills.sqlite; a writer/reader path-split
# would leave the endpoint empty (the ict-health-snapshot artifacts-path class
# of bug).
#
# Save to: /etc/systemd/system/ict-exchange-fills-pull.service
# (installed automatically by scripts/install_systemd_units.sh during
#  deploy_pull_restart.sh)
# ============================================

[Unit]
Description=ICT exchange-fills pull — keep /api/bot/pnl/exchange fresh (BL-20260713-EXCHANGE-FILLS-STORE-EMPTY)
After=network-online.target ict-trader-live.service
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/ict-trading-bot
# Load the trader .env for the Bybit creds (BYBIT_API_KEY_2 / _SECRET_2).
EnvironmentFile=-/home/ubuntu/ict-trading-bot/.env
# Pin DATA_DIR so runtime_state_dir() resolves to the block-volume root the
# web-api reads from, regardless of whether .env carries it.
Environment=DATA_DIR=/data/bot-data
# Quick Bybit API call (~5s for a day's fills); cap generously for a slow
# session. A read/auth failure exits non-zero and the unit is marked failed —
# the next timer fire retries; no state is left half-written (exec_id PK).
TimeoutStartSec=300
# Run under the trader venv so ccxt imports (same interpreter the trader uses
# for Bybit) — not system /usr/bin/python3.
ExecStart=/home/ubuntu/ict-trading-bot/.venv/bin/python -u scripts/pull_exchange_fills.py --account bybit_2 --category linear --days 7 --api-key-env BYBIT_API_KEY_2 --api-secret-env BYBIT_API_SECRET_2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
