# deploy/ict-trader-live.service
[Unit]
Description=ICT Trader Live
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/ict-trading-bot
EnvironmentFile=/home/ubuntu/ict-trading-bot/.env
# Resource priority (2026-06-10, trader-connection-stability). The trader's
# single-threaded main loop + the heartbeat it writes between ticks must NEVER
# be starved by a misbehaving sidecar — the 2026-06-10 incident was exactly
# that: the IB-gateway watchdog's docker-restart churn (heavy Java) saturated
# the 2-core VM (loadavg ~10), the trader loop couldn't get scheduled, the
# heartbeat froze, and the liveness watchdog (also starved) couldn't autoheal
# for ~25 min. High CPUWeight + IOWeight (cgroup-v2 shares; translated to
# CPUShares on v1) + a negative Nice guarantee the trader CPU + I/O first under
# contention. These only bind UNDER contention — an idle box is unaffected.
CPUWeight=10000
IOWeight=10000
Nice=-10
# 2026-05-08 operator directive: tick every minute (still on 5-min
# candles — strategy timeframe unchanged). HEARTBEAT_INTERVAL stays
# at 60 s so liveness signal fires in lock-step with the tick. Both
# can be overridden in /home/ubuntu/ict-trading-bot/.env if a future
# operator wants a slower cadence; the values here are unit-level
# defaults that win when the .env doesn't set them.
Environment=TICK_INTERVAL_SECONDS=60
Environment=HEARTBEAT_INTERVAL_SECONDS=60
# Enable the DB→Bybit reconciler (src/runtime/order_monitor.py::
# _reconcile_open_trades) on every monitor tick. Without this flag
# the reconciler is gated off and the trade journal can drift from
# exchange truth — e.g. trade 1145 on 2026-05-11 stayed open in DB
# while the Bybit position had already closed. Pair with the reverse
# reconciler (_reconcile_orphan_exchange_positions) which runs in the
# same tick and ADOPTs any Bybit position that has no journal row.
Environment=MONITOR_RECONCILE_ENABLED=true
# Reverse-reconciler policy when a Bybit position has no matching
# open trade row. detect_only = alert + log, no mutation. adopt =
# insert a new trade row so the monitor manages it from the next
# tick. close = submit a market close. See
# src/runtime/order_monitor.py::ORPHAN_POSITION_POLICY for semantics.
Environment=ORPHAN_POSITION_POLICY=adopt
# Multi-strategy intent layer (PR #1125 + #1129). When true the
# pipeline runs the intent-aware multiplexer
# (src/runtime/intent_multiplexer.py): every enabled strategy emits
# a StrategyIntent each tick, the aggregator collapses them into one
# DesiredPosition for BTCUSDT (same-direction reinforcement keeps the
# larger valid target; conflicts resolve deterministically by
# priority), and Coordinator.multi_account_execute sends only the
# delta to reach that target. Per-account RiskManager still owns
# sizing and the daily-loss / margin / position-size caps —
# the intent layer never bypasses risk.
# Default is false in src/runtime/intent_multiplexer.py; pinning to
# true here is the operator's opt-in for Bybit2 / BTCUSDT (Turtle
# Soup + VWAP coexistence). Override or revert by editing this line
# and firing pull-and-deploy.
Environment=MULTI_STRATEGY_INTENT_LAYER=true
ExecStart=/usr/bin/python3 -u -B -m src.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
