# deploy/ict-ib-gateway-watchdog.timer
# ============================================
# ICT IB GATEWAY WATCHDOG TIMER (reactive auto-restart re-armed 2026-06-22)
# Fires ict-ib-gateway-watchdog.service every ~5 min.
#
# History: 2026-06-10 the gateway-isolation redesign dropped this to once-daily
# alert-only, leaning on the deterministic 05:30 `docker restart`
# (ict-ib-gateway-reset.timer) for recovery. That left a real gap — a session
# that wedges MID-DAY (between 05:30 resets) had no automatic recovery and
# stayed dark until the next day (observed 2026-06-22,
# BL-20260622-GATEWAY-MIDDAY-WEDGE). Now that the gateway is on its own VM, the
# reactive restart churn can't touch the money loop, so the bounded reactive
# probe is re-armed (guards live in the .service ExecStart). ~5 min cadence ⇒ a
# mid-day wedge is caught within ~10 min (--restart-after 2). The daily
# ict-ib-gateway-reset.timer stays as a belt-and-suspenders scheduled restart
# (retimed to 06:05 UTC 2026-07-02, was 05:30 — BL-20260623-002; the .service
# ExecStart now also carries --suppress-window-utc 03:45-05:45 so this watchdog
# never burns a restart attempt on a wedge it can't fix during IBKR's own reset).
#
# Auto-installed + enabled by scripts/install_systemd_units.sh (gateway VM only).
# Save to: /etc/systemd/system/ict-ib-gateway-watchdog.timer
# ============================================

[Unit]
Description=ICT IB Gateway watchdog — reactive probe + bounded auto-restart (~5 min)
Requires=ict-ib-gateway-watchdog.service

[Timer]
# Reactive cadence: first probe 3 min after boot, then every 5 min. Bounded
# auto-restart guards in the .service keep this loop-proof.
OnBootSec=3min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=ict-ib-gateway-watchdog.service

[Install]
WantedBy=timers.target
