# deploy/ict-prop-telegram-listener.service
# ============================================
# ICT PROP TELEGRAM LISTENER
# Inbound half of the Breakout manual bridge: long-polls the prop bot's
# getUpdates and turns an operator's structured command (e.g.
# "close ETHUSD 2950 +80 tp") into a /api/bot/prop/report ingest — no
# Claude/dashboard middle-man. See src/prop/telegram_inbound.py.
#
# Activates only when a prop bot token AND an allowed chat id are set
# (TELEGRAM_PROP_ALLOWED_CHAT_IDS / TELEGRAM_CHAT_ID); otherwise it logs
# and exits (Restart=always will retry, harmless until configured).
#
# Save to: /etc/systemd/system/ict-prop-telegram-listener.service
# ============================================

[Unit]
Description=ICT Trading Bot - Prop Telegram Listener (inbound report-back)
After=network-online.target ict-trader-live.service
Wants=network-online.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/ict-trading-bot
ExecStart=/usr/bin/python3 -u -B -m src.prop.telegram_inbound
Restart=always
RestartSec=15
# Nice/low-weight: a long-poll loop, must never contend with the trader.
Nice=10
StandardOutput=journal
StandardError=journal
EnvironmentFile=/home/ubuntu/ict-trading-bot/.env

[Install]
WantedBy=multi-user.target
