# deploy/ict-bot.service
[Unit]
Description=ICT Trading Bot
After=network.target

[Service]
# Run as the ubuntu user that owns the repo
User=ubuntu
Group=ubuntu

# Absolute paths for the Oracle Cloud Ubuntu instance
WorkingDirectory=/home/ubuntu/ict-trading-bot
EnvironmentFile=/home/ubuntu/ict-trading-bot/.env
Environment=PYTHONPATH=/home/ubuntu/ict-trading-bot

# Launch via module so relative imports resolve correctly
ExecStart=/usr/bin/python3 -m src.main

# Restart policy: restart on crash, not on clean exit
Restart=on-failure
RestartSec=30

# Log to systemd journal (view with: journalctl -u ict-bot -f)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
