# deploy/ict-git-sync.service
# ============================================
# ICT GIT SYNC SYSTEMD SERVICE
# Pulls latest code and restarts services when triggered
# Save to: /etc/systemd/system/ict-git-sync.service
# ============================================

[Unit]
Description=ICT Git Auto Update
After=network.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/ict-trading-bot
ExecStart=/bin/bash /home/ubuntu/ict-trading-bot/scripts/deploy_pull_restart.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
