# deploy/ict-ib-gateway-reset.timer
# ============================================
# ICT IB GATEWAY DAILY RESET TIMER (gateway-isolation redesign, 2026-06-10)
# Fires ict-ib-gateway-reset.service once daily at 05:30 UTC — just after
# IBKR's overnight US server-reset window (~03:45–05:45 UTC) and after the
# image's own IBC AUTO_RESTART_TIME (23:59 UTC), so the full container restart
# lands when IBKR is healthy again and re-establishes a clean session for the
# trading day. One deterministic restart replaces the reactive 5-min restart
# loop. Tune the time here if IBKR's reset window shifts.
#
# Save to: /etc/systemd/system/ict-ib-gateway-reset.timer
# ============================================

[Unit]
Description=ICT IB Gateway daily reset — fires at 05:30 UTC
Requires=ict-ib-gateway-reset.service

[Timer]
OnCalendar=*-*-* 05:30:00 UTC
# If the box was off at 05:30, run on next boot so a reset is never skipped.
Persistent=true
AccuracySec=1min
Unit=ict-ib-gateway-reset.service

[Install]
WantedBy=timers.target
