# deploy/ict-ib-gateway-watchdog.timer
# ============================================
# ICT IB GATEWAY HEALTH-CHECK TIMER (gateway-isolation redesign, 2026-06-10)
# Fires ict-ib-gateway-watchdog.service once DAILY (alert-only).
#
# Was: every 5 min, reactive auto-restart loop. Now: the deterministic daily
# `docker restart` on the gateway VM (ict-ib-gateway-reset.timer) is the
# recovery; this is just a once-daily alert-only probe so a persistently-dark
# MES still pages the operator. Fires at 06:00 UTC — 30 min after the daily
# reset (05:30 UTC) so it checks the post-reset session. MES is delayed paper
# data, so daily-granularity alerting is fine; the trader-side connect breaker
# handles intraday blips without paging.
#
# Auto-installed + enabled by scripts/install_systemd_units.sh.
# Save to: /etc/systemd/system/ict-ib-gateway-watchdog.timer
# ============================================

[Unit]
Description=ICT IB Gateway health-check — fires daily (alert-only)
Requires=ict-ib-gateway-watchdog.service

[Timer]
# Once daily at 06:00 UTC — 30 min after the daily reset (05:30) so it probes
# the freshly-reset session.
OnCalendar=*-*-* 06:00:00 UTC
AccuracySec=1min
Persistent=true
Unit=ict-ib-gateway-watchdog.service

[Install]
WantedBy=timers.target
