# deploy/ict-ib-gateway-watchdog.service
# ============================================
# ICT IB GATEWAY WATCHDOG SYSTEMD SERVICE
# Runs scripts/check_ib_gateway.py once per timer fire.
#
# Auto-heals the recurring MES/IB Gateway wedge (BL-20260527-003): the
# ib-gateway container stays "up" but its IBKR session dies during IBKR's
# overnight server reset (data farms broken, every account/market-data
# request times out), so ib_paper collects no MES candles and all MES
# strategies skip every tick. The fix IBC can't self-apply is a full
# container restart — this watchdog probes ib_paper and, after a sustained
# wedge, runs scripts/ops/restart_ib_gateway.sh (the same docker-restart
# the manual vm-ib-gateway-recover dispatches).
#
# Safety: --restart-after 2 (need ~10 min of wedge before acting),
# --max-restarts 3 + --cooldown-min 20 so a genuinely bad credential /
# IBKR lockout can NEVER turn into a restart loop (every restart is a fresh
# IBKR login). Once exhausted the watchdog alert-only escalates — but
# --exhaustion-reset-min 120 re-arms that budget after 2h of no restart, so a
# multi-hour wedge whose first 3 restarts all land inside IBKR's overnight
# server-reset window (when a re-login can't yet re-establish the session) is
# retried once IBKR recovers, instead of stranding MES for the whole episode
# (the 2026-06-09 incident, BL-20260605-004). Still loop-proof: the re-arm
# only fires after a long quiet gap.
#
# Requires passwordless sudo for `docker restart` (same as the
# vm-ib-gateway-recover workflow path) and ib_insync importable under the
# system python3 (same interpreter the trader uses).
#
# Save to: /etc/systemd/system/ict-ib-gateway-watchdog.service
# (installed automatically by scripts/install_systemd_units.sh during deploy)
# ============================================

[Unit]
Description=ICT Trading Bot — IB Gateway Watchdog (MES auto-heal)
After=network-online.target ict-trader-live.service
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/ict-trading-bot
# `-` prefix: TOLERATE the file being absent (BL-20260622-GATEWAY-NO-ENV).
# This unit runs on the dedicated GATEWAY VM, which is provisioned with only
# the IB credential env + the Docker container — it has NO repo-root `.env`.
# Without the `-`, systemd hard-fails the start ("Failed to load environment
# files: No such file or directory" → result 'resources'), so the re-armed
# watchdog could never run at all (verified 2026-06-22). The .env only carries
# the Telegram creds for send_alert(), which already degrades gracefully when
# unset; the probe + the docker-restart self-heal need no env. Same `-` pattern
# as ict-git-sync.service.
EnvironmentFile=-/home/ubuntu/ict-trading-bot/.env
# NOTE: state path is pinned to the repo's runtime_logs/ via --state below, NOT
# $DATA_DIR. The gateway VM has no writable /data/bot-data (the canonical OCI
# data dir exists only on the trader), so the old DATA_DIR=/data/bot-data made
# save_state fail with "Permission denied: '/data'" → the wedged streak never
# persisted, so the watchdog could never reach --restart-after (verified
# 2026-06-22, BL-20260622-GATEWAY-LOCAL-PROBE). runtime_logs/ is repo-local +
# ubuntu-writable and survives `git reset --hard` (untracked).
# Resource priority (2026-06-10): so the probe + any `docker restart` it
# dispatches run promptly even under load. The IB-gateway container itself is
# CPU/memory-capped at `docker run` (scripts/install_ib_gateway_docker.sh:
# --cpus/--memory), and `docker restart` preserves those flags — so the
# restart churn this watchdog can produce during an IBKR reset window can no
# longer saturate the 2-core VM (the 2026-06-10 cascade). The watchdog stays
# high-priority so it itself isn't starved while acting.
CPUWeight=9000
Nice=-10
# REACTIVE AUTO-RESTART (re-armed 2026-06-22, BL-20260622-GATEWAY-MIDDAY-WEDGE).
# History: the reactive probe+restart was DISARMED on 2026-06-10 (gateway-
# isolation redesign) in favour of one daily deterministic `docker restart`
# (ict-ib-gateway-reset.timer), because the reactive restart churn used to
# saturate the 2-core box the gateway SHARED with the trader. That objection is
# now moot — the gateway runs on its OWN dedicated VM, so its restart churn can
# never touch the money loop. Meanwhile the daily-only posture left a real gap:
# a session that wedges MID-DAY (between 05:30 resets) had NO automatic recovery
# and stayed dark until the next 05:30 — observed 2026-06-22 (the gateway died
# ~07:3x UTC after the 05:30 reset; MES/MGC/MHG were dark and the open MHG
# position pkg-f58a249d119049ab tripped a MONITOR BLIND alert; only a manual
# vm-ib-gateway-recover dispatch restored it). Re-arming the SAME bounded guards
# closes that gap while staying loop-proof:
#   --restart-after 2      need ~10 min of sustained wedge before acting
#   --max-restarts 3       cap per episode (a bad credential / IBKR lockout can
#                          never become a restart loop — every restart is a
#                          fresh IBKR login)
#   --cooldown-min 20      min gap between restarts
#   --exhaustion-reset-min 120  re-arm the budget after 2h quiet, so a multi-hour
#                          wedge whose first restarts land inside IBKR's reset
#                          window is retried once IBKR recovers (BL-20260605-004)
# Once the budget is exhausted the watchdog falls back to alert-only until the
# session recovers. The manual vm-ib-gateway-recover workflow remains as the
# operator escape hatch. The reactive cadence lives in the sibling .timer.
#
# DEP-FREE LOCAL PROBE (--probe-script, BL-20260622-GATEWAY-LOCAL-PROBE): the
# gateway VM is a MINIMAL box with no bot venv, so the account probe
# (ib_connect_check.py) can't import ib_insync and falsely reads
# "connect failed" (→ would loop-restart a healthy gateway). ib_gateway_local_probe.py
# diagnoses the wedge from the container's own state + recent `docker logs`
# (socat→127.0.0.1:4002 refused / IBC re-auth pending, no recent
# "Login has completed") using ONLY the docker CLI — no venv, no creds. --state
# is pinned to the repo-local (writable) runtime_logs/ so the wedged streak
# actually persists between the timer's 5-min runs.
# SUPPRESSION WINDOW (2026-07-02, BL-20260623-002): 03:45-05:45 UTC matches
# IBKR's own documented overnight reset window (see
# deploy/ict-ib-gateway-reset.timer) — a wedge detected in here can't be
# fixed by a restart, so this narrowly stops the watchdog from burning its
# restart budget/cooldown on an attempt that was never going to work
# (that's what pushed the effective restart past the observed 06:00-06:05Z
# recurrence). Still logs/alerts on detection for visibility; the streak is
# frozen, not reset, so it resumes the instant the window closes. Scoped to
# exactly this window — the mid-day auto-heal re-armed 2026-06-22
# (BL-20260622-GATEWAY-MIDDAY-WEDGE) is unaffected outside it.
ExecStart=/usr/bin/python3 -u scripts/check_ib_gateway.py --probe-account ib_paper \
  --probe-script scripts/ops/ib_gateway_local_probe.py \
  --state runtime_logs/ib_gateway_watchdog_state.json \
  --auto-restart --restart-after 2 --max-restarts 3 \
  --cooldown-min 20 --exhaustion-reset-min 120 \
  --suppress-window-utc 03:45-05:45
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
