# deploy/ict-shadow-log-rotate.service
# ICT shadow-prediction audit log rotation (S-AI-WS7-FU).
#
# ENABLED + ACTIVE on the trader VM (field beats this comment's old "disabled
# by default" wording — corrected 2026-06-28 full-system audit). The installer's
# auto-enable loop (scripts/install_systemd_units.sh) runs `enable --now` on every
# deploy/*.timer that is NOT in _GATEWAY_ONLY_TIMERS — and this timer is not — so
# it runs on the trader box with no manual opt-in. Live-verified `enabled active`
# via system-action probe #4946; added to diag _CANONICAL_UNITS in #4947. (A
# manual `sudo systemctl enable --now ict-shadow-log-rotate.timer` still applies
# on a host where the installer never ran.)
#
# Rotates runtime_logs/shadow_predictions.jsonl when it exceeds
# size or age thresholds. Always exits 0 — log rotation must not
# crash the system. Errors are journaled.
#
# Runs as the same user that writes the log (ubuntu, by convention).
# Reads the working tree at /opt/ict-trading-bot/ (symlink to the
# canonical clone). Idempotent — re-runs are safe.
[Unit]
Description=ICT shadow-prediction audit log rotation
After=local-fs.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/opt/ict-trading-bot
ExecStart=/usr/bin/python3 scripts/ops/rotate_shadow_log.py --gzip
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
