# deploy/ict-web-api.service
# systemd unit — ICT Trading Bot dashboard API.
#
# Save to: /etc/systemd/system/ict-web-api.service
#
# S-014: binds to 0.0.0.0:8001 so the Vercel React dashboard can reach
# the API from the browser. Restrict inbound access at the firewall
# (ufw allow 8001/tcp) if needed.
#
# DASHBOARD_ORIGIN is set here directly so the git-sync autosync picks
# it up without a manual edit on the VM. The /api/bot/* endpoints are
# read-only telemetry; no secret data is exposed.
#
# Required env (from /etc/ict-trader/web-api.env):
#   JWT_SIGNING_KEY    — HS256 signing key (for /api/auth/* routes)
#   ALLOWED_EMAIL      — single-operator allowlist
#   GIT_SHA            — optional, used by runtime_status when git CLI
#                        is unavailable inside the unit
#
# To apply after autosync:
#   sudo systemctl daemon-reload
#   sudo systemctl restart ict-web-api
#
# To enable on a fresh VM:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now ict-web-api

[Unit]
Description=ICT Trading Bot — dashboard API
After=network-online.target ict-trader-live.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/ict-trading-bot
Environment=PYTHONPATH=/opt/ict-trading-bot
Environment=DASHBOARD_ORIGIN=https://bentzbk-ict-trader-dashboard.vercel.app
EnvironmentFile=/etc/ict-trader/web-api.env
ExecStart=/usr/bin/python3 -m uvicorn src.web.api.main:app --host 0.0.0.0 --port 8001
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
