# deploy/trainer/ict-promotion-readiness.service
# ============================================
# ICT PROMOTION-READINESS SWEEP (S-MLOPT-S18, M14 Phase 4.3)
# Runs scripts/ops/run_promotion_readiness.sh once per fire.
#
# Sweeps the whole registry via `python -m ml promotion-readiness` and
# writes report.json + SUMMARY.md under
#   runtime_logs/trainer_mirror/promotion_readiness/<UTC-date>/
# so the existing publish_trainer_mirror.sh rsync picks it up. On any
# actionable proposal (promote-ready / demote-proposed), drops a
# pending_pings/*.json onto the live VM via SCP so
# ict-telegram-bot.service delivers the operator ping.
#
# **TRAINER VM ONLY** — never the Oracle live trader. Lives in
# deploy/trainer/ so the live-VM unit installer
# (scripts/install_systemd_units.sh, which globs deploy/*.service)
# never picks it up.
#
# Save to (trainer VM): /etc/systemd/system/ict-promotion-readiness.service
# Install:  cp deploy/trainer/ict-promotion-readiness.{service,timer} /etc/systemd/system/
#           && systemctl daemon-reload
#           && systemctl enable --now ict-promotion-readiness.timer
#
# Mirrors the unit embedded in deploy/training-vm-cloud-init.yaml so a
# freshly provisioned trainer + a long-running trainer end up with the
# same unit body. Edits must stay in lock-step.
# ============================================

[Unit]
Description=ICT promotion-readiness sweep (S-MLOPT-S18)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/ict-trading-bot
Environment=LIVE_VM_IP=141.145.193.91
Environment=LIVE_VM_USER=ubuntu
Environment=VM_SSH_KEY=/home/ubuntu/.ssh/id_ed25519
# 2026-07-19 (WS-A / readiness-starvation): the default 1h heavy-lock queue
# wait starved the sweep on BOTH 07-18 and 07-19 — the timer fired ~04:10Z,
# waited its full 3600s, and gave up while the nightly training cycle (which
# has grown past a 3h runtime) still held the queue. The queue's design is
# "wait your turn", so give the sweep a 4h wait: it then runs right after the
# cycle releases the lock instead of silently skipping the day. (The 04:00
# retime of 2026-07-17 fixed the CONCURRENT-memory collision; this fixes the
# serialized-starvation half.)
Environment=TRAINER_HEAVY_LOCK_WAIT_S=14400
# Hard memory bound (MB-20260719-PROMOREADY-OOSEDGE-OOM): the sweep
# D-state-thrashed the whole 6 GB box at ~5 GB RSS twice on 2026-07-19 with
# nothing to stop it. A breach now OOM-kills the service fast (unit fails,
# timer retries next day) instead of hours of page-thrash. MemorySwapMax=0
# keeps a breach from degenerating back into thrash via swap.
MemoryMax=4500M
MemorySwapMax=0
ExecStart=/bin/bash scripts/ops/run_promotion_readiness.sh
# Exit code 10 means "sweep ran + actionable proposals found" and is a
# normal outcome (the operator ping was queued) — do not let systemd
# treat it as a service failure.
SuccessExitStatus=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
